# deploy/supervisor-idsea.conf
[program:idsea-backend]
command=/var/www/idsea/backend/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8001 --workers 2 --log-level info
directory=/var/www/idsea/backend
user=www-data
autostart=true
autorestart=true
startretries=5
stopwaitsecs=30
stderr_logfile=/var/log/idsea/backend.err.log
stdout_logfile=/var/log/idsea/backend.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
environment=PATH="/var/www/idsea/backend/venv/bin:%(ENV_PATH)s"
